# rootfs/etc/systemd/system/gem.service
[Service]
Type=oneshot
ExecStart=bash /opt/gem/u/gem.sh